# powersave/powersave.service
[Unit]
Description=powersave

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/powersave

[Install]
WantedBy=multi-user.target
